# deploy/systemd/fix-my-claw-oneshot.service
[Unit]
Description=fix-my-claw OpenClaw watchdog (oneshot)
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot

# Run as the same user that runs OpenClaw (recommended).
# User=openclaw

Environment=PYTHONUNBUFFERED=1
# Render @FIX_MY_CLAW_BIN@ and @CONFIG_PATH@ before installing.
ExecStart=@FIX_MY_CLAW_BIN@ repair --config @CONFIG_PATH@
