# deploy/systemd/fix-my-claw.service
[Unit]
Description=fix-my-claw OpenClaw watchdog
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=always
RestartSec=5

# Run as the same user that runs OpenClaw (recommended).
# User=openclaw

Environment=PYTHONUNBUFFERED=1
# Render @FIX_MY_CLAW_BIN@ and @CONFIG_PATH@ before installing.
ExecStart=@FIX_MY_CLAW_BIN@ monitor --config @CONFIG_PATH@

[Install]
WantedBy=multi-user.target
